# varnish/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/varnishd
drwxr-xr-x	root/root	etc/varnish/
-rw-r--r--	root/root	etc/varnish/default.vcl
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/varnishadm
-rwxr-xr-x	root/root	usr/bin/varnishhist
-rwxr-xr-x	root/root	usr/bin/varnishlog
-rwxr-xr-x	root/root	usr/bin/varnishncsa
-rwxr-xr-x	root/root	usr/bin/varnishstat
-rwxr-xr-x	root/root	usr/bin/varnishtest
-rwxr-xr-x	root/root	usr/bin/varnishtop
drwxr-xr-x	root/root	usr/include/
drwxr-xr-x	root/root	usr/include/varnish/
drwxr-xr-x	root/root	usr/include/varnish/cache/
-rw-r--r--	root/root	usr/include/varnish/cache/cache.h
-rw-r--r--	root/root	usr/include/varnish/cache/cache_backend.h
-rw-r--r--	root/root	usr/include/varnish/cache/cache_director.h
-rw-r--r--	root/root	usr/include/varnish/cache/cache_filter.h
-rw-r--r--	root/root	usr/include/varnish/cache/cache_varnishd.h
drwxr-xr-x	root/root	usr/include/varnish/common/
-rw-r--r--	root/root	usr/include/varnish/common/common_param.h
-rw-r--r--	root/root	usr/include/varnish/miniobj.h
drwxr-xr-x	root/root	usr/include/varnish/tbl/
-rw-r--r--	root/root	usr/include/varnish/tbl/acct_fields_bereq.h
-rw-r--r--	root/root	usr/include/varnish/tbl/acct_fields_req.h
-rw-r--r--	root/root	usr/include/varnish/tbl/backend_poll.h
-rw-r--r--	root/root	usr/include/varnish/tbl/ban_vars.h
-rw-r--r--	root/root	usr/include/varnish/tbl/bo_flags.h
-rw-r--r--	root/root	usr/include/varnish/tbl/boc_state.h
-rw-r--r--	root/root	usr/include/varnish/tbl/body_status.h
-rw-r--r--	root/root	usr/include/varnish/tbl/cli_cmds.h
-rw-r--r--	root/root	usr/include/varnish/tbl/debug_bits.h
-rw-r--r--	root/root	usr/include/varnish/tbl/feature_bits.h
-rw-r--r--	root/root	usr/include/varnish/tbl/h2_error.h
-rw-r--r--	root/root	usr/include/varnish/tbl/h2_frames.h
-rw-r--r--	root/root	usr/include/varnish/tbl/h2_settings.h
-rw-r--r--	root/root	usr/include/varnish/tbl/h2_stream.h
-rw-r--r--	root/root	usr/include/varnish/tbl/htc.h
-rw-r--r--	root/root	usr/include/varnish/tbl/http_headers.h
-rw-r--r--	root/root	usr/include/varnish/tbl/http_response.h
-rw-r--r--	root/root	usr/include/varnish/tbl/locks.h
-rw-r--r--	root/root	usr/include/varnish/tbl/obj_attr.h
-rw-r--r--	root/root	usr/include/varnish/tbl/oc_exp_flags.h
-rw-r--r--	root/root	usr/include/varnish/tbl/oc_flags.h
-rw-r--r--	root/root	usr/include/varnish/tbl/params.h
-rw-r--r--	root/root	usr/include/varnish/tbl/req_body.h
-rw-r--r--	root/root	usr/include/varnish/tbl/req_flags.h
-rw-r--r--	root/root	usr/include/varnish/tbl/sess_attr.h
-rw-r--r--	root/root	usr/include/varnish/tbl/sess_close.h
-rw-r--r--	root/root	usr/include/varnish/tbl/steps.h
-rw-r--r--	root/root	usr/include/varnish/tbl/symbol_kind.h
-rw-r--r--	root/root	usr/include/varnish/tbl/vcc_types.h
-rw-r--r--	root/root	usr/include/varnish/tbl/vcl_returns.h
-rw-r--r--	root/root	usr/include/varnish/tbl/vhd_fsm.h
-rw-r--r--	root/root	usr/include/varnish/tbl/vhd_fsm_funcs.h
-rw-r--r--	root/root	usr/include/varnish/tbl/vhd_return.h
-rw-r--r--	root/root	usr/include/varnish/tbl/vhp_huffman.h
-rw-r--r--	root/root	usr/include/varnish/tbl/vhp_static.h
-rw-r--r--	root/root	usr/include/varnish/tbl/vrt_stv_var.h
-rw-r--r--	root/root	usr/include/varnish/tbl/vsc_levels.h
-rw-r--r--	root/root	usr/include/varnish/tbl/vsl_tags.h
-rw-r--r--	root/root	usr/include/varnish/tbl/vsl_tags_http.h
-rw-r--r--	root/root	usr/include/varnish/tbl/waiters.h
drwxr-xr-x	root/root	usr/include/varnish/vapi/
-rw-r--r--	root/root	usr/include/varnish/vapi/vapi_options.h
-rw-r--r--	root/root	usr/include/varnish/vapi/voptget.h
-rw-r--r--	root/root	usr/include/varnish/vapi/vsc.h
-rw-r--r--	root/root	usr/include/varnish/vapi/vsl.h
-rw-r--r--	root/root	usr/include/varnish/vapi/vsl_int.h
-rw-r--r--	root/root	usr/include/varnish/vapi/vsm.h
-rw-r--r--	root/root	usr/include/varnish/vas.h
-rw-r--r--	root/root	usr/include/varnish/vav.h
-rw-r--r--	root/root	usr/include/varnish/vbm.h
-rw-r--r--	root/root	usr/include/varnish/vcl.h
-rw-r--r--	root/root	usr/include/varnish/vcli.h
-rw-r--r--	root/root	usr/include/varnish/vcs.h
-rw-r--r--	root/root	usr/include/varnish/vdef.h
-rw-r--r--	root/root	usr/include/varnish/vmod_abi.h
-rw-r--r--	root/root	usr/include/varnish/vqueue.h
-rw-r--r--	root/root	usr/include/varnish/vre.h
-rw-r--r--	root/root	usr/include/varnish/vrnd.h
-rw-r--r--	root/root	usr/include/varnish/vrt.h
-rw-r--r--	root/root	usr/include/varnish/vrt_obj.h
-rw-r--r--	root/root	usr/include/varnish/vsa.h
-rw-r--r--	root/root	usr/include/varnish/vsb.h
-rw-r--r--	root/root	usr/include/varnish/vsha256.h
-rw-r--r--	root/root	usr/include/varnish/vtcp.h
-rw-r--r--	root/root	usr/include/varnish/vtim.h
-rw-r--r--	root/root	usr/include/varnish/vut.h
-rw-r--r--	root/root	usr/include/varnish/vut_options.h
drwxr-xr-x	root/root	usr/include/varnish/waiter/
-rw-r--r--	root/root	usr/include/varnish/waiter/waiter.h
drwxr-xr-x	root/root	usr/lib/
-rwxr-xr-x	root/root	usr/lib/libvarnishapi.la
lrwxrwxrwx	root/root	usr/lib/libvarnishapi.so -> libvarnishapi.so.2.0.0
lrwxrwxrwx	root/root	usr/lib/libvarnishapi.so.2 -> libvarnishapi.so.2.0.0
-rwxr-xr-x	root/root	usr/lib/libvarnishapi.so.2.0.0
drwxr-xr-x	root/root	usr/lib/pkgconfig/
-rw-r--r--	root/root	usr/lib/pkgconfig/varnishapi.pc
drwxr-xr-x	root/root	usr/lib/varnish/
drwxr-xr-x	root/root	usr/lib/varnish/vmods/
-rwxr-xr-x	root/root	usr/lib/varnish/vmods/libvmod_blob.la
-rwxr-xr-x	root/root	usr/lib/varnish/vmods/libvmod_blob.so
-rwxr-xr-x	root/root	usr/lib/varnish/vmods/libvmod_directors.la
-rwxr-xr-x	root/root	usr/lib/varnish/vmods/libvmod_directors.so
-rwxr-xr-x	root/root	usr/lib/varnish/vmods/libvmod_proxy.la
-rwxr-xr-x	root/root	usr/lib/varnish/vmods/libvmod_proxy.so
-rwxr-xr-x	root/root	usr/lib/varnish/vmods/libvmod_purge.la
-rwxr-xr-x	root/root	usr/lib/varnish/vmods/libvmod_purge.so
-rwxr-xr-x	root/root	usr/lib/varnish/vmods/libvmod_std.la
-rwxr-xr-x	root/root	usr/lib/varnish/vmods/libvmod_std.so
-rwxr-xr-x	root/root	usr/lib/varnish/vmods/libvmod_unix.la
-rwxr-xr-x	root/root	usr/lib/varnish/vmods/libvmod_unix.so
-rwxr-xr-x	root/root	usr/lib/varnish/vmods/libvmod_vtc.la
-rwxr-xr-x	root/root	usr/lib/varnish/vmods/libvmod_vtc.so
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/varnishd
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/aclocal/
-rw-r--r--	root/root	usr/share/aclocal/varnish-legacy.m4
-rw-r--r--	root/root	usr/share/aclocal/varnish.m4
drwxr-xr-x	root/root	usr/share/doc/
drwxr-xr-x	root/root	usr/share/doc/varnish/
-rw-r--r--	root/root	usr/share/doc/varnish/builtin.vcl
-rw-r--r--	root/root	usr/share/doc/varnish/example.vcl
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man1/
-rw-r--r--	root/root	usr/share/man/man1/varnishadm.1.gz
-rw-r--r--	root/root	usr/share/man/man1/varnishd.1.gz
-rw-r--r--	root/root	usr/share/man/man1/varnishhist.1.gz
-rw-r--r--	root/root	usr/share/man/man1/varnishlog.1.gz
-rw-r--r--	root/root	usr/share/man/man1/varnishncsa.1.gz
-rw-r--r--	root/root	usr/share/man/man1/varnishstat.1.gz
-rw-r--r--	root/root	usr/share/man/man1/varnishtest.1.gz
-rw-r--r--	root/root	usr/share/man/man1/varnishtop.1.gz
drwxr-xr-x	root/root	usr/share/man/man3/
-rw-r--r--	root/root	usr/share/man/man3/vmod_blob.3.gz
-rw-r--r--	root/root	usr/share/man/man3/vmod_directors.3.gz
-rw-r--r--	root/root	usr/share/man/man3/vmod_proxy.3.gz
-rw-r--r--	root/root	usr/share/man/man3/vmod_purge.3.gz
-rw-r--r--	root/root	usr/share/man/man3/vmod_std.3.gz
-rw-r--r--	root/root	usr/share/man/man3/vmod_unix.3.gz
-rw-r--r--	root/root	usr/share/man/man3/vmod_vtc.3.gz
drwxr-xr-x	root/root	usr/share/man/man7/
-rw-r--r--	root/root	usr/share/man/man7/varnish-cli.7.gz
-rw-r--r--	root/root	usr/share/man/man7/varnish-counters.7.gz
-rw-r--r--	root/root	usr/share/man/man7/vcl.7.gz
-rw-r--r--	root/root	usr/share/man/man7/vsl-query.7.gz
-rw-r--r--	root/root	usr/share/man/man7/vsl.7.gz
-rw-r--r--	root/root	usr/share/man/man7/vtc.7.gz
drwxr-xr-x	root/root	usr/share/varnish/
drwxr-xr-x	root/root	usr/share/varnish/vcl/
-rw-r--r--	root/root	usr/share/varnish/vcl/devicedetect.vcl
-rwxr-xr-x	root/root	usr/share/varnish/vmodtool.py
-rwxr-xr-x	root/root	usr/share/varnish/vsctool.py
drwxr-xr-x	root/root	usr/var/
drwxr-xr-x	root/root	usr/var/varnish/
